# proton-pass-ssh-agent.service
[Unit]
Description=Proton Pass SSH Agent (pass-cli native)
Documentation=https://github.com/dhruvesh16/proton_sshagent
After=graphical-session.target
Wants=graphical-session.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=%h/.local/bin/proton-pass-ssh-agent-wrapper.sh
Restart=always
RestartSec=5
Environment=DISPLAY=:0

# Security hardening
ProtectHome=false
NoNewPrivileges=true

[Install]
WantedBy=default.target
